# templates/admin/etc/cron.d/swift_reports
# By default it will send mail to the swiftops alias below
# which in turn redirects it to the proper location 
##MAILTO="swiftops"
## m h  dom mon dow   command
55 23 * * 0 root /usr/local/bin/swift_recon_wrapper.sh
##58 23 * * * swiftops bash /usr/local/bin/storage_usage.sh
